# playbooks/roles/webserver/templates/supervisor_microservice.conf.j2
[program:{{ microservice_name }}]
command=npm start
directory={{ microservice_home }}
user=www-data 